# cmd/jib-watcher/jib-watcher.service
[Unit]
Description=Jib Git Watcher
After=jib-deployer.service
Wants=jib-deployer.service

[Service]
ExecStartPre=/bin/sh -c 'until nc -z 127.0.0.1 4222; do sleep 1; done'
ExecStart=/usr/local/bin/jib-watcher run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
