# cmd/jib-bus/jib-bus.service
[Unit]
Description=Jib Bus (NATS via docker compose)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c 'docker compose -f /opt/jib/stack/docker-compose.yml -p jib-stack up -d jib-bus'
ExecStop=/bin/sh -c 'docker compose -f /opt/jib/stack/docker-compose.yml -p jib-stack stop jib-bus'

[Install]
WantedBy=multi-user.target
