# cmd/jib-cloudflared/jib-cloudflared.service
[Unit]
Description=Jib Cloudflared Tunnel (via docker compose)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c 'docker compose -f /opt/jib/cloudflared/docker-compose.yml up -d'
ExecStop=/bin/sh -c 'docker compose -f /opt/jib/cloudflared/docker-compose.yml down'

[Install]
WantedBy=multi-user.target
